# docker-nginx-flask/deployment/application.conf
description "uWSGI server instance configured to serve flask application"

start on runlevel [2345]
stop on runlevel [!2345]

setuid root
setgid www-data

chdir /srv
exec uwsgi --ini nginx/uwsgi.ini
